# systemd/dynadash-backend.service
[Unit]
Description=DynaDash FastAPI Backend
After=network.target influxdb.service
Wants=influxdb.service

[Service]
Type=simple
User=root
WorkingDirectory=/home/user/dynadash/backend
ExecStart=/home/user/dynadash/backend/.venv/bin/uvicorn main:app \
    --host 127.0.0.1 \
    --port 8000 \
    --log-level info
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=dynadash-backend

# Allow the backend to write its data directory
ReadWritePaths=/home/user/dynadash/backend/data

[Install]
WantedBy=multi-user.target
